# 6/config/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/etc/init.d/mysqld start

[program:apache]
command=/usr/sbin/apache2ctl -D FOREGROUND
